# docker/resources/supervisord.conf
[supervisord]
nodaemon=true
user = root

[program:nginx]
user = root
command=/usr/sbin/nginx -g "daemon off;"

[program:php-fpm]
user = root
command=/usr/sbin/php-fpm7.4 -F